# Qmail/etc/service/mysqld/run
#!/bin/sh
exec /sbin/setuser mysql mysqld 2>&1
